# smart_boiler/rootfs/etc/services.d/smart_boiler/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start Smart Boiler service
# ==============================================================================
exec python3 -m smart-boiler.control --url "http://supervisor/core/api" --key "$SUPERVISOR_TOKEN" --addon "True"