# smartboiler/rootfs/etc/services.d/smartboiler/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start SMART BOILER service
# ==============================================================================

exec python3 -m smartboiler.control --url "http://supervisor/core/api" --key "$SUPERVISOR_TOKEN" --addon "True"
